# build_files/components/core/auto-updates/files/usr/lib/systemd/system/rpm-ostree-countme.service.d/10-override.conf
# Inert while rpm-ostree-countme.timer is masked (enable-services.sh).
# If the timer is unmasked, this makes the ping retry instead of leaving a
# failed unit when it runs before the network is up, e.g. right after
# resume from suspend.
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
Restart=on-failure
RestartSec=5m
