# build_files/files/common/usr/lib/systemd/user/seed-justfile.service
[Unit]
Description=Seed ~/.config/just/justfile for existing users
ConditionPathExists=!%h/.seed-justfile-done

[Service]
Type=oneshot
ExecStart=/usr/libexec/seed-justfile
ExecStartPost=/usr/bin/touch %h/.seed-justfile-done

[Install]
WantedBy=default.target
