# build_files/files/common/usr/lib/systemd/user/update-notify.service
[Unit]
Description=Notify the user of a staged system update

[Service]
Type=oneshot
ExecStart=/usr/libexec/update-notify
