# build_files/files/common/usr/lib/systemd/system/bootc-fetch-apply-updates.service.d/10-override.conf
[Unit]
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Restart=on-failure
RestartSec=60s
# Stock unit runs `bootc upgrade --apply`, which reboots immediately once an
# update is fetched. Drop --apply so it only stages the update; update-notify
# (user session) then prompts to restart when ready.
ExecStart=
ExecStart=/usr/bin/bootc upgrade --quiet
